# mem0-archive.timer
[Unit]
Description=mem0 Memory Archive Timer
Documentation=https://github.com/norrishuang/mem0-memory-service
Requires=mem0-archive.service

[Timer]
# Run daily at 02:00 UTC (10:00 Beijing Time)
OnCalendar=*-*-* 02:00:00
Persistent=true
# If system was offline, run on next boot
OnBootSec=5min

[Install]
WantedBy=timers.target
