# systemd/mem0-memory.service
[Unit]
Description=mem0 Memory Service - OpenClaw Agent Memory Layer
After=network.target

[Service]
Type=simple
User=REPLACE_WITH_YOUR_USER
WorkingDirectory=REPLACE_WITH_INSTALL_PATH
EnvironmentFile=REPLACE_WITH_INSTALL_PATH/.env
ExecStart=/usr/bin/python3 REPLACE_WITH_INSTALL_PATH/server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
